# build/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:xvfb]
command=/usr/bin/Xvfb %(ENV_DISPLAY)s -screen 0 %(ENV_RESOLUTION)sx24 -ac -pn -noreset +extension GLX +extension RANDR +extension RENDER -nolisten tcp
priority=100
user=%(ENV_USER)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/xvfb.log
stderr_logfile=/var/log/supervisor/xvfb.err
environment=DISPLAY="%(ENV_DISPLAY)s"

[program:pulseaudio]
command=/usr/bin/pulseaudio --disallow-module-loading --disallow-exit --exit-idle-time=-1
priority=200
user=%(ENV_USER)s
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/pulseaudio.log
stderr_logfile=/var/log/supervisor/pulseaudio.err

[program:desktop]
command=/usr/bin/xfce4-session
priority=300
user=%(ENV_USER)s
autostart=true
autorestart=false
stdout_logfile=/var/log/supervisor/desktop.log
stderr_logfile=/var/log/supervisor/desktop.err
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s"

[program:selkies]
command=selkies-gstreamer --port=8081
directory=/home/%(ENV_USER)s
priority=400
user=%(ENV_USER)s
autostart=true
autorestart=true
stdout_logfile=/tmp/selkies.log
stderr_logfile=/tmp/selkies.log
environment=DISPLAY="%(ENV_DISPLAY)s",HOME="%(ENV_HOME)s",USER="%(ENV_USER)s"